# sdk-webapp.service
[Unit]
Description=sdk webapp

[Service]
ExecStart=/usr/lib/sdk-webapp-bundle/sdk-webapp.sh 
